# apt/lvm-snapshot/etc/systemd/system/lvm-snapshot.timer
[Unit]
Description=Daily LVM snapshot for root volume

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

